# files/etc/uci-defaults/99-custom
uci -q batch << EOI
set network.lan.proto='dhcp'
delete network.lan.ipaddr
set mjpg-streamer.core.enabled='1'
EOI
